# scripts/com.agentic.trading-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: the SLOW PLANNER pass every 15 minutes (StartInterval 900), around the clock. This
  is the deliberative half of the two-rate loop (docs/two-rate-architecture.md): discovery + DD +
  manage + entry commit/arm (the LLM work). The FAST half — protective exits + armed-entry triggers,
  no LLM — runs every minute via com.agentic.sentinel.plist. The two share data/.tick.lock, so they
  never race; a DD-heavy planner tick just makes that minute's sentinel pass skip. 15 min gives a DD
  tick room to finish (a DD is ~85s) and matches the multi-day drift horizon; the once-a-day DD cache
  keeps token cost bounded. (Was 600 = 10 min.)

  Off-hours ticks are NOT failures. market_conditions.py + tick_context.py compute the session in
  America/New_York and print GATE=SKIP outside weekday 09:30-15:55 ET; run_trading_tick.sh then logs
  a deterministic skip and exits 0 with NO LLM call. So the market-hours filtering lives in the
  script (the correct backstop), and launchd just needs a plain 1-min cadence — no need to enumerate
  hundreds of StartCalendarInterval slots (it has no cron-style step syntax anyway).

  Cost of a clean cadence: a cheap SKIP log line every 1 min outside market hours. Trivial.

  Install:  cp scripts/com.agentic.trading-tick.plist ~/Library/LaunchAgents/
            launchctl unload -w ~/Library/LaunchAgents/com.agentic.trading-tick.plist 2>/dev/null
            launchctl load   -w ~/Library/LaunchAgents/com.agentic.trading-tick.plist
  Stop:     launchctl unload -w ~/Library/LaunchAgents/com.agentic.trading-tick.plist
  Note: the Mac must be AWAKE for ticks to fire (launchd won't wake it from sleep; it will run
  one coalesced tick on wake if a scheduled time was missed).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic.trading-tick</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mcosta/agentic-trading/scripts/run_trading_tick.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mcosta/agentic-trading</string>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/mcosta/agentic-trading/data/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/mcosta/agentic-trading/data/logs/launchd.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
