# scripts/com.agentic.paper-sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: PAPER sentinel fast pass every 1 minute. Deterministic protective exits +
  armed-entry triggers, NO LLM. Pairs with com.agentic.trading-paper.plist (10-min planner).

  Install:  cp scripts/com.agentic.paper-sentinel.plist ~/Library/LaunchAgents/
            launchctl unload -w ~/Library/LaunchAgents/com.agentic.paper-sentinel.plist 2>/dev/null
            launchctl load   -w ~/Library/LaunchAgents/com.agentic.paper-sentinel.plist
  Stop:     launchctl unload -w ~/Library/LaunchAgents/com.agentic.paper-sentinel.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic.paper-sentinel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/ABSOLUTE/PATH/TO/agentic-trading/scripts/run_paper_sentinel.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/agentic-trading</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/agentic-trading/data/logs/paper-sentinel.out</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/agentic-trading/data/logs/paper-sentinel.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
